# deploy/pomment-next.service
[Unit]
Description=Pomment Next
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pomment
Group=pomment
WorkingDirectory=/opt/pomment-next
Environment=NODE_ENV=production
EnvironmentFile=/etc/pomment-next.env
ExecStart=/usr/local/bin/bun run start
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

StateDirectory=pomment-next
StateDirectoryMode=0700
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/pomment-next

[Install]
WantedBy=multi-user.target
